# scripts/s1-time-display.service
[Unit]
Description=AceMagic S1 Time Display Service
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/s1-display/src/time_display.py
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
